# ops/supervisord.conf
[supervisord]
nodaemon=false

[program:app]
directory=/app/src
command=/bin/bash -c "npm start"

[program:gulp]
directory=/app/src
command=/bin/bash -c "gulp"